# systemd/elephant-backup-systemd.timer
[Unit]
Description=Elephant Backup's auto snapshots.

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
