# systemd/elephant-backup-systemd.service
[Unit]
Description=Elephant Backup's auto snapshots.
After=zfs.target

[Service]
Type=simple
ExecStart="elephant-backup snapshot POOLS_REPLACEMENT"
