# systemd/elephant-backup-snapshot.service
[Unit]
After=zfs.target
[Service]
Type=simple
ExecStart=elephant-backup snapshot POOLS_REPLACEMENT
